# internal/apps/catalog/manifests/filebrowser.yaml
id: filebrowser
name: File Browser
description: Web file manager
category: tools
web_path: /filebrowser/
detection:
  binaries:
    - filebrowser
  systemd_units:
    - filebrowser.service
