# internal/apps/catalog/manifests/plex.yaml
id: plex
name: Plex
description: Media server
category: media
web_path: /plex/
detection:
  systemd_units:
    - plexmediaserver.service
