# contrib/systemd/brrewery.service
[Unit]
Description=brrewery package management API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/brrewery serve
Restart=on-failure
RestartSec=5
# brrewery compiles and installs system packages through Ansible. When it runs
# as a non-root user it escalates per operation with sudo, using a password
# entered in the web UI. Sandboxing such as NoNewPrivileges= (blocks sudo) or
# ProtectSystem=strict (makes /usr, /etc, /opt read-only even for the privileged
# child) would prevent package installation, so they are intentionally not set.

[Install]
WantedBy=multi-user.target
