# internal/packages/catalog/manifests/emby.yaml
id: emby
name: Emby
description: Media server
category: media
web_path: /emby/
detection:
  systemd_units:
    - emby-server.service
