# internal/packages/catalog/manifests/lidarr.yaml
id: lidarr
name: Lidarr
description: Music management
category: arr
web_path: /lidarr/
detection:
  binaries:
    - lidarr
  systemd_units:
    - lidarr.service
