# internal/packages/catalog/manifests/radarr.yaml
id: radarr
name: Radarr
description: Movie management
category: arr
web_path: /radarr/
detection:
  binaries:
    - radarr
  systemd_units:
    - radarr.service
