# internal/packages/catalog/manifests/sonarr.yaml
id: sonarr
name: Sonarr
description: TV series management
category: arr
web_path: /sonarr/
detection:
  binaries:
    - sonarr
  systemd_units:
    - sonarr.service
