# ansible/roles/deluge_build/templates/deluge-web@.service.j2
[Unit]
Description=Deluge WebUI for %i
After=network-online.target deluged@%i.service
Wants=network-online.target
Requires=deluged@%i.service

[Service]
Type=simple
User=%i
{% if dz_legacy_ld | length > 0 %}
Environment=LD_LIBRARY_PATH={{ dz_legacy_ld }}
{% endif %}
# dz_web_fg is -d on the 2.x lines (do-not-daemonize); 1.3's deluge-web already
# runs in the foreground and has no such flag. nginx reverse-proxies this port
# under {{ deluge_web_base }}, so bind it to localhost and tell deluge-web its
# base path (-b) so it emits correct asset URLs behind the proxy. No -c: it uses
# ~/.config/deluge via the HOME systemd sets from User=%i (not %h, which is /root
# in a system unit).
ExecStart={{ deluge_web_bin }} {{ dz_web_fg }} -i 127.0.0.1 -p {{ deluge_web_port }} -b {{ deluge_web_base }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
